# trunk/install/adminconsole.conf
# adminconsole - Admin web console for gproxy

description "admin console"

start on runlevel [2345]
stop on runlevel [!2345]

#expect fork
respawn

chdir /www/admin

exec /www/admin/app.py
